# cookbooks/github_action/files/home/github-action/crontab
0 3 * * * /usr/local/github_action/import/rake import:daily
